# confs/crontab
1 0 * * * /vagrant/scripts/nightly.sh
0 0 1 * * /vagrant/scripts/monthly.sh
5 8 * * 6 /vagrant/scripts/weekly.sh
