# confs/crontab
1 0 * * * root /vagrant/scripts/nightly.sh
0 0 1 * * root /vagrant/scripts/monthly.sh
5 8 * * 6 root /vagrant/scripts/weekly.sh
